# resync/config/resync.service
# =============================================================================
# Systemd Service File for RESYNC v6.2.0
# =============================================================================
#
# Installation:
#   sudo cp config/resync.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable resync
#   sudo systemctl start resync
#
# Management:
#   Status:   sudo systemctl status resync
#   Logs:     sudo journalctl -u resync -f
#   Restart:  sudo systemctl restart resync
#   Stop:     sudo systemctl stop resync
#
# Production command (same as ExecStart below):
#   gunicorn -c gunicorn_config.py resync.main:app
#
# =============================================================================

[Unit]
Description=Resync v6.2.0 - AI-Powered TWS/HWA Operations Platform
Documentation=https://github.com/resync/resync
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=notify
# Gunicorn supports systemd readiness notifications (sd_notify).
# See Gunicorn systemd deployment docs.
NotifyAccess=main

# User and Group
User=resync
Group=resync

# Working Directory
WorkingDirectory=/opt/resync

# Environment
EnvironmentFile=/opt/resync/.env
Environment="PATH=/opt/resync/.venv/bin:/usr/local/bin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="ENVIRONMENT=production"

# ── Production Start Command ─────────────────────────────────────────────────
# Gunicorn master + Uvicorn async workers
# Config: gunicorn_config.py (workers, timeouts, logging)
ExecStart=/opt/resync/.venv/bin/gunicorn \
    -c /opt/resync/gunicorn_config.py \
    resync.main:app

# ── Restart Policy ───────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitInterval=60s

# ── Timeouts ─────────────────────────────────────────────────────────────────
# Allow time for Docling model loading on first request
TimeoutStartSec=180
TimeoutStopSec=30

# ── Resource Limits ──────────────────────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=4096

# Adjust per server capacity:
# MemoryMax=8G
# MemoryHigh=6G

# ── Security Hardening ───────────────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/resync /var/log/resync /tmp

# Capabilities
AmbientCapabilities=
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
